# .deploy_settings/supervisor/supervisor-uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/app/bin/uwsgi -i /srv/app/.deploy_settings/uwsgi/shippa_order.ini